# mkosi.conf.d/11-koji/mkosi.extra/usr/lib/systemd/system/kojid.service.d/add-host.conf
[Service]
ExecStartPre=koji add-host localhost.localdomain x86-64
ExecStartPre=koji add-host-to-channel localhost.localdomain createrepo
ExecStartPre=koji edit-host --capacity=4 localhost.localdomain