# mkosi.conf.d/11-koji/mkosi.extra/usr/lib/systemd/system/httpd.service.d/add-tag.conf
[Service]
ExecStartPost=koji add-tag mkosi
ExecStartPost=koji add-tag --parent mkosi --arches x86_64 mkosi-build
ExecStartPost=koji add-external-repo -t mkosi-build mkosi-repo https://mirrors.kernel.org/fedora/development/39/Everything/x86_64/os
ExecStartPost=koji add-target mkosi mkosi-build
ExecStartPost=koji add-group mkosi-build build
ExecStartPost=koji add-group-pkg mkosi-build build mkosi
